# deploy/rashodomer.service
[Unit]
Description=Rashodomer FastAPI Backend
After=network.target

[Service]
WorkingDirectory=/var/rashodomer-be/src
ExecStart=/var/rashodomer-be/.venv/bin/uvicorn app:fast_api --host 0.0.0.0 --port 80
Restart=always
User=root
Environment="PATH=/var/rashodomer-be/.venv/bin"

[Install]
WantedBy=multi-user.target
